# sync/com.kernelspace.obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kernelspace.obsidian-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/edwinabrahamthomas/edwinat/kernel-space-|/sync/sync-to-obsidian.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/edwinabrahamthomas/edwinat/kernel-space-|/sync/launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/edwinabrahamthomas/edwinat/kernel-space-|/sync/launchd-error.log</string>
</dict>
</plist>
